# archive-stream-00.service
[Unit]
Description=Archive fcso stream at top of the hour

[Service]
Type=oneshot
ExecStart=%h/bin/archive-stream.sh

